# etc/services.d/recoll-webui/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Start Recoll Web UI
# ==============================================================================

exec /opt/recollwebui/webui-standalone.py -a 0.0.0.0
